# systemd/collectord.service
# jz_sniff_rn — collectord: Event dedup, SQLite persistence, JSON export
#
# Depends on sniffd (receives events via IPC).  Writes to /var/lib/jz/jz.db.

[Unit]
Description=jz_sniff_rn Collector Daemon
Documentation=file:/usr/share/doc/jz_sniff_rn/design.md
After=sniffd.service
BindsTo=sniffd.service
PartOf=sniffd.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=forking
ExecStart=/usr/local/sbin/collectord --daemon --config /etc/jz/base.yaml
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/jz/collectord.pid

Restart=on-failure
RestartSec=5

LimitNOFILE=65536

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/jz /run/jz
ProtectHome=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=collectord

[Install]
WantedBy=multi-user.target
