# systemd/sniffd.service
# jz_sniff_rn — sniffd: BPF loader, event consumer, probe generator, REST API
#
# Must start after rSwitch (provides XDP pipeline) and network.
# Requires CAP_BPF / CAP_NET_ADMIN; LimitMEMLOCK=infinity for BPF maps.

[Unit]
Description=jz_sniff_rn Sniff Daemon
Documentation=file:/usr/share/doc/jz_sniff_rn/design.md
After=network-online.target rswitch.service
Wants=network-online.target rswitch.service configd.service collectord.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=forking
ExecStartPre=/bin/sh -c 'rm -rf /sys/fs/bpf/jz && rm -f /sys/fs/bpf/jz_* /sys/fs/bpf/rs_ctx_map /sys/fs/bpf/rs_progs /sys/fs/bpf/rs_prog_chain /sys/fs/bpf/rs_event_bus && mkdir -p /sys/fs/bpf/jz'
ExecStart=/usr/local/sbin/sniffd --daemon --config /etc/jz/base.yaml
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/jz/sniffd.pid

Restart=on-failure
RestartSec=5

# BPF maps require unlimited locked memory
LimitMEMLOCK=infinity
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/jz /run/jz /sys/fs/bpf/jz /sys/fs/bpf /etc/netplan
ProtectHome=yes
PrivateTmp=yes

# Runtime directory (created by systemd, cleaned on stop)
RuntimeDirectory=jz
RuntimeDirectoryPreserve=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sniffd

[Install]
WantedBy=multi-user.target
